# archlinux/etc/systemd/system/systemd-fsck@.service.d/override.conf
[Service]
StandardOutput=null
StandardError=journal+console
TimeoutSec=
TimeoutSec=0
