# archlinux/etc/systemd/system/ly.service.d/override.conf
[Unit]
After=
After=systemd-user-sessions.service plymouth-quit-wait.service
After=getty@tty1.service

[Service]
TTYPath=
TTYPath=/dev/tty1
